# solaris-build/proto/lib/svc/manifest/vmware/balloon.xml
<?xml version="1.0" ?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<!--
 WMware::Ballooning

    Service manifest to declare and install the service that
    will namage the memory monitor for VMware balloon
-->

<service_bundle type='manifest' name='vmware:balloon'>
    <service version='1' type='service' name='vmware/balloon'>

        <!-- Instance -->

        <create_default_instance enabled='false'/>
        <single_instance/>

        <!-- Dependencies -->

        <dependency restart_on="none" type="service"
            name="multi_user_dependency" grouping="require_all">
            <service_fmri value="svc:/milestone/multi-user"/>
        </dependency>

        <dependency restart_on="error" type="service"
            name="filesystem_dependency" grouping="require_all">
            <service_fmri value='svc:/system/filesystem/minimal' />
        </dependency>

        <!-- Methods -->

        <exec_method timeout_seconds='30' type='method' name='start'
            exec='/lib/svc/method/vmware-balloon %m'/>

        <exec_method timeout_seconds='30' type='method' name='restart'
            exec='/lib/svc/method/vmware-balloon %m'/>

        <exec_method timeout_seconds='30' type='method' name='stop'
            exec='/lib/svc/method/vmware-balloon %m'/>

        <exec_method timeout_seconds='30' type='method' name='refresh'
            exec=':true'/>

        <!-- Properties -->

        <property_group name='startd' type='framework'>
                <propval name='duration' type='astring' value='transient' />
        </property_group>

        <stability value='Unstable' />

        <!-- Description -->

        <template>
            <common_name>
                <loctext xml:lang='C'>
                    Memory balloon control service for VMware
                </loctext>
            </common_name>
            <description>
                <loctext xml:lang='C'>
                    Will enable the memory balloon support if requested by the hypervisor.
                </loctext>
            </description>
        </template>
    </service>
</service_bundle>
