# solaris-build/proto/lib/svc/manifest/vmware/tools.xml
<?xml version="1.0" ?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<!--
 VMware::WMware-Tools

    Service manifest to declare and install the service that
    will manage the vmware daemon to monitor and administrate
    the virtual machine from the hypervisor
-->

<service_bundle type='manifest' name='vmware:tools'>
    <service version='1' type='service' name='vmware/tools'>

        <!-- Dependencies -->

        <dependency restart_on="none" type="service"
            name="multi_user_dependency" grouping="require_all">
            <service_fmri value="svc:/milestone/multi-user"/>
        </dependency>

        <dependency restart_on="error" type="service"
            name="filesystem_dependency" grouping="require_all">
            <service_fmri value='svc:/system/filesystem/minimal' />
        </dependency>

        <!-- Methods -->

        <exec_method timeout_seconds='30' type='method' name='start'
            exec='/lib/svc/method/vmware-tools %m'/>

        <exec_method timeout_seconds='30' type='method' name='restart'
            exec='/lib/svc/method/vmware-tools %m'/>

        <exec_method timeout_seconds='30' type='method' name='stop'
            exec=':kill'/>

        <exec_method timeout_seconds='30' type='method' name='refresh'
            exec=':true'/>

        <!-- Instance -->

        <instance enabled="true" name="default"/>
        <stability value='Unstable' />

        <!-- Description -->

        <template>
            <common_name>
                <loctext xml:lang='C'>
                    Open VM Tools for VMware
                </loctext>
            </common_name>
            <description>
                <loctext xml:lang='C'>
                    Will start the vmware vmtoolsd daemon to monitor and communicate with the VM from the hypervisor.
                </loctext>
            </description>
        </template>
    </service>
</service_bundle>
